# scripts/board-watchdog.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!-- board-watchdog.plist.template — the launchd LaunchAgent for the #1435 external
     board-freshness watchdog. NOT a live plist: scripts/install-board-watchdog.sh
     substitutes the __PLACEHOLDER__ tokens with real absolute paths at install time
     (so this committed template carries NO home path — the CI privacy gate stays
     green). It mirrors the com.user.cc-mailbox-doorbell periodic StartInterval house
     pattern (the daily jobs use StartCalendarInterval instead). -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>__LABEL__</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd '__REPO__' &amp;&amp; npm run --silent kanban:watchdog</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>__HOME__</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>StartInterval</key>
  <integer>__INTERVAL__</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>WorkingDirectory</key>
  <string>__REPO__</string>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/kanban-board-watchdog.out.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/kanban-board-watchdog.err.log</string>
</dict>
</plist>
